# external/openwrt-prebuilts/x86/uci-defaults/0_default_config
#!/bin/sh
opkg install /ipks/*
ip link set eth0 mtu 1460
