# packages/modules/Bluetooth/system/gd/rust/chromeos/examples/upstart/btmanagerd.conf
description     "Start the bluetooth manager daemon"
author          "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services

exec btmanagerd